# asl3-herald.service
[Unit]
Description=ASL3 Herald - Enhanced Tail Message & Announcement Daemon
Documentation=https://github.com/N6LKA/asl3-herald
After=network.target asterisk.service
Wants=asterisk.service

[Service]
Type=simple
ExecStart=/usr/local/bin/asl3-herald/asl3-herald.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asl3-herald

# Reduce unnecessary restarts during Asterisk maintenance
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
